# scrolllock-monitor.service
[Unit]
Description=Scroll Lock Monitor Service
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/scrolllock_monitor.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
